# docker/config/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * for i in {1..6}; do php -e /app/tools/mailer.php; sleep 10; done
